# scripts/com.agent-runtime.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-runtime.backup</string>

    <!-- Daily meta-repo backup at 02:00. macOS launchd alternative to the cron
         entry in install-cron.sh. Rendered + loaded by
         scripts/install-timers-macos.sh. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{REPO_ROOT}}/scripts/backup-wiki.sh</string>
    </array>

    <!-- Daily at 02:00. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>2</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{REPO_ROOT}}/log/cron-backup.log</string>
    <key>StandardErrorPath</key>
    <string>{{REPO_ROOT}}/log/cron-backup.log</string>

    <key>WorkingDirectory</key>
    <string>{{REPO_ROOT}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- backup-wiki needs git + a PyYAML python (it probes) to read
             backup.meta_remote. Mirror the service plist PATH. -->
        <key>PATH</key><string>/opt/homebrew/bin:{{HOME}}/bin:{{HOME}}/.local/bin:{{HOME}}/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NO_PROXY</key><string>*</string>
        <key>no_proxy</key><string>*</string>
        <key>http_proxy</key><string></string>
        <key>https_proxy</key><string></string>
        <key>HTTP_PROXY</key><string></string>
        <key>HTTPS_PROXY</key><string></string>
    </dict>
</dict>
</plist>
